# .config/systemd/user/waybar.service.d/10-resource-limits.conf
[Service]
# 2026-07-06 Waybar tray path reached 14.8G and triggered global OOM.
# Keep tray enabled, but contain future leaks inside the user service cgroup.
MemoryHigh=768M
MemoryMax=1G
Restart=on-failure
RestartSec=2s
