# .config/systemd/user/plasma-xembedsniproxy.service.d/10-resource-limits.conf
[Service]
# Bridge legacy XEmbed tray icons into StatusNotifierItem for Waybar.
# Keep it bounded because it participates in the tray path.
MemoryHigh=192M
MemoryMax=256M
Restart=on-failure
RestartSec=2s
