# .config/systemd/user/mako.service.d/10-user-config.conf
[Service]
# Let systemd restart the notification daemon if it exits unexpectedly.
Restart=on-failure
RestartSec=2s
