# .config/systemd/user/cliphist.service.d/10-user-config.conf
[Service]
# Keep the previous niri startup behavior while letting systemd own the process.
ExecStart=
ExecStart=/usr/bin/wl-paste --watch cliphist store -max-items 5000
Restart=on-failure
RestartSec=2s
